# assets/runit/wss/finish
#!/bin/bash
killall websocket-server.js
