# assets/runit/wss/run
#!/bin/bash
/var/www/h/npm/websocket-server.js
