# assets/runit/wss/run
#!/bin/bash
/usr/bin/nodejs /var/www/h/npm/websocket-server.js
