# assets/runit/cool/run
#!/bin/bash
python /var/www/h/assets/cooling/cool.py
